# roles/nut/templates/nut_exporter.service.j2
# {{ ansible_managed }}
# =============================================================================
# NUT Exporter Systemd Service
# =============================================================================
# Prometheus exporter for Network UPS Tools (NUT)
# https://github.com/DRuggeri/nut_exporter

[Unit]
Description=NUT Exporter for Prometheus
Documentation=https://github.com/DRuggeri/nut_exporter
After=network-online.target nut-server.service
Wants=network-online.target

[Service]
Type=simple
User={{ nut_exporter_user }}
Group={{ nut_exporter_group }}

ExecStart={{ nut_exporter_install_dir }}/nut_exporter \
  --nut.server={{ nut_exporter_nut_server }} \
  --nut.serverport={{ nut_exporter_nut_port }} \
{% if nut_exporter_nut_username %}
  --nut.username={{ nut_exporter_nut_username }} \
{% endif %}
  --web.listen-address={{ nut_exporter_listen_address }}:{{ nut_exporter_port }} \
  --web.telemetry-path={{ nut_exporter_metrics_path }} \
{% if nut_exporter_variables | length > 0 %}
  --nut.vars_enable={{ nut_exporter_variables | join(',') }} \
{% endif %}
  --log.level={{ nut_exporter_log_level }}

{% if nut_exporter_nut_password %}
Environment=NUT_EXPORTER_PASSWORD={{ nut_exporter_nut_password }}
{% endif %}

Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
